# monitoring/scrubber/scrubber.service
[Unit]
Description=BTRFS scrub service

[Service]
Type=oneshot
EnvironmentFile=/etc/default/scrubber
ExecStart=scrubber ${DEVICES} -c ${CONTACT}